# meta-ecofleet/recipes-ecofleet/gobi-ui/files/gobi-ui.service
[Unit]
Description=EcoFleet Gobi Dashboard
After=graphical.target gobi-agent.service
Wants=graphical.target

[Service]
Type=simple
User=root
Environment=QT_QPA_PLATFORM=wayland
Environment=WAYLAND_DISPLAY=wayland-1
Environment=XDG_RUNTIME_DIR=/run/user/0
# Wait for Weston socket before starting
ExecStartPre=/bin/sh -c 'i=0; while [ ! -S "${XDG_RUNTIME_DIR}/wayland-1" ] && [ $i -lt 30 ]; do sleep 1; i=$((i+1)); done'
ExecStart=/usr/sbin/gobi-ui
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
